# systemd/sacaoos.service
[Unit]
Description=SecAIOS Telegram Bot with OpenCode
After=network.target

[Service]
Type=simple
User=aios
WorkingDirectory=/opt/sacaoos
ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/python3 execution/telegram/opencode_bot.py
Restart=always
RestartSec=10

Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target